# power_button/power-button.service
[Unit]
Description=Raspberry Pi Power Button Service
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python3 /home/pi/power_button.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/usr/local/lib/python3.9/site-packages

[Install]
WantedBy=multi-user.target
